# packages/matrix-rssbot/matrix_rssbot-0.1.0.tar.gz/matrix_rssbot-0.1.0/rssbot-pantalaimon.service
[Unit]
Description=Pantalaimon for RSSbot
Requires=network.target

[Service]
Type=simple
User=rssbot
Group=rssbot
WorkingDirectory=/opt/rssbot
ExecStart=/opt/rssbot/venv/bin/python3 -um pantalaimon.main -c pantalaimon.conf
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target